# dashy.service
[Unit]
Description=dashy service control plane
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/server/dashy
ExecStart=python3 /home/ubuntu/server/dashy/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
